# game/RedstoneForge/Client/rootfs/etc/supervisor/conf.d/bot.conf
[program:bot]
command=/waitjava.sh bot bot.jar
redirect_stderr=true
autostart=true
stdout_syslog=true
user=minecraft

